# config/supervisor/brand-search-worker.conf
[program:brand-search-worker]
command=php /home/eskill/htdocs/eskill.com.br/bin/brand-search-worker.php --daemon
directory=/home/eskill/htdocs/eskill.com.br
user=eskill
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=30
stderr_logfile=/home/eskill/htdocs/eskill.com.br/storage/logs/brand-search-worker.err.log
stdout_logfile=/home/eskill/htdocs/eskill.com.br/storage/logs/brand-search-worker.log
stderr_logfile_maxbytes=10MB
stdout_logfile_maxbytes=10MB
